# src/fujin/templates/default.timer.j2
[Unit]
Description=Timer for {{ app_name }} {{ process_name }}

[Timer]
OnCalendar={{ process.timer }}
Unit={{ process_name }}.service

[Install]
WantedBy=timers.target
